# scripts/systemd/sw-cortex-slack.service
[Unit]
Description=sw-cortex Slack Handler (Socket Mode)
After=network.target

[Service]
Type=simple
User=jackk
WorkingDirectory=/home/jackk/sw-cortex
Environment="PATH=/home/jackk/.nvm/versions/node/v22.17.0/bin:/usr/bin"
ExecStart=/home/jackk/.nvm/versions/node/v22.17.0/bin/npx tsx scripts/slack-handler.ts
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
